# workers/news-ingest/deploy/systemd/situation-monitor-news-worker.service
[Unit]
Description=Situation Monitor News Ingest Worker
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/situation-monitor-news-worker
ExecStart=/usr/bin/bash /opt/situation-monitor-news-worker/run.sh once
